# config/supervisor/galaxy/gunicorn.ini
[program:galaxy-gunicorn]
command         = gunicorn 'galaxy.webapps.galaxy.fast_factory:factory()' --timeout 300 --pythonpath lib -k galaxy.webapps.galaxy.workers.Worker -b %(ENV_GALAXY_BIND_IP)s:%(ENV_GALAXY_PORT)s
directory       = %(ENV_GALAXY_ROOT_DIR)s
process_name    = gunicorn
umask           = %(ENV_GALAXY_UMASK)s
autostart       = true
autorestart     = true
startsecs       = 15
stopwaitsecs    = 65
environment     = UNICORNHERDER_LOGLEVEL=info
numprocs        = 1
stdout_logfile  = %(ENV_GALAXY_LOG_DIR)s/gunicorn.log
redirect_stderr = true
