# config/supervisor/galaxy/unicornherder.ini
[program:galaxy-unicornherder]
command         = unicornherder --pidfile %(ENV_GALAXY_MUTABLE_DATA_DIR)s/herder.pid -- 'galaxy.webapps.galaxy.fast_factory:factory()' --timeout 300 --pythonpath lib -k galaxy.webapps.galaxy.workers.Worker -b %(ENV_GALAXY_BIND_IP)s:%(ENV_GALAXY_PORT)s --access-logfile %(ENV_GALAXY_LOG_DIR)s/gunicorn.access.log --error-logfile %(ENV_GALAXY_LOG_DIR)s/gunicorn.error.log --capture-output
directory       = %(ENV_GALAXY_ROOT_DIR)s
process_name    = unicornherder
umask           = 022
autostart       = true
autorestart     = true
startsecs       = 15
stopwaitsecs    = 65
user            = %(ENV_GALAXY_USER)s
environment     = UNICORNHERDER_LOGLEVEL=info
numprocs        = 1
stdout_logfile  = %(ENV_GALAXY_LOG_DIR)s/unicornherder.log
redirect_stderr = true
